# deploy/oracle/preview-runner/systemd/preview-runner.service
[Unit]
Description=AutoDocker Preview Runner
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=autodocker
Group=autodocker
WorkingDirectory=/srv/preview-runner/app
EnvironmentFile=/srv/preview-runner/env/preview-runner.env
ExecStart=/srv/preview-runner/app/.venv/bin/gunicorn config.wsgi:application --bind 127.0.0.1:9000 --workers 2 --timeout 120
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
